# packaging/openonedrived.service.in
[Unit]
Description=open-onedrive custom FUSE OneDrive daemon
After=default.target

[Service]
Type=simple
ExecStart=@INSTALL_BIN_DIR@/openonedrived
Restart=on-failure
RestartSec=3
Environment=RUST_LOG=openonedrived=info

[Install]
WantedBy=default.target
